# webserver/python_wsgi/supervisor_app.conf
[program:mongo-geoipdb]
user = user
command=/path/to/virtualenv/bin/uwsgi /path/project/uwsgi/uwsgi.ini
autostart=true
autorestart=true
stopsignal = QUIT
